# supervisord/gpsclock.ini
[program:gpsclock]
command=/opt/gpsclock/main.py
autorestart=true
startsecs=0
user=root
redirect_stderr=true
